# supervisord.conf
[supervisord]
nodaemon = true
logfile = /dev/stdout
loglevel = info

[unix_http_server]
file = /run/supervisord.s

[supervisorctl]
serverurl = unix:///run/supervisord.s

[include]
files = /etc/supervisord/conf.d/*.conf

[program:icinga2]
command = /bin/bash -exo pipefail -c '. /usr/lib/icinga2/icinga2; . /etc/default/icinga2; exec icinga2 daemon -e "$ICINGA2_ERROR_LOG"'
stdout_logfile = /dev/stdout
redirect_stderr = true
autostart = true
autorestart = true
startsecs = 10
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 5

[program:apache2]
command = /bin/bash -exo pipefail -c '. /etc/apache2/envvars; exec apache2 -k start -DFOREGROUND'
stdout_logfile = /dev/stdout
redirect_stderr = true
autostart = true
autorestart = true
startsecs = 10
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 5

[program:mariadb]
command = /usr/sbin/mysqld -u mysql
stdout_logfile = /dev/stdout
redirect_stderr = true
autostart = true
autorestart = true
startsecs = 10
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 5

[program:influxdb]
command = /bin/bash -exo pipefail -c '. /etc/default/influxdb; exec nosu influxdb influxdb influxd -config /etc/influxdb/influxdb.conf $INFLUXD_OPTS'
stdout_logfile = /dev/stdout
redirect_stderr = true
autostart = true
autorestart = true
startsecs = 10
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 5

[program:grafana]
command = /bin/bash -exo pipefail -c 'cd /usr/share/grafana; . /etc/default/grafana-server; exec nosu grafana grafana grafana-server "--config=$CONF_FILE" "--pidfile=/var/run/grafana/grafana-server.pid" "cfg:default.paths.logs=$LOG_DIR" "cfg:default.paths.data=$DATA_DIR" "cfg:default.paths.plugins=$PLUGINS_DIR" "cfg:default.paths.provisioning=$PROVISIONING_CFG_DIR"'
stdout_logfile = /dev/stdout
redirect_stderr = true
autostart = true
autorestart = true
startsecs = 10
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 5
